# deploy/printable.service
[Unit]
Description=Printable Box (POS receipt dispenser)
After=network-online.target pigpiod.service
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APP_DIR__
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=__APP_DIR__/.env
# Neustart bei Absturz UND nach einem Self-Update (Prozess beendet sich mit 0)
Restart=always
RestartSec=3
# Zugriff auf GPIO/Audio benötigt ggf. diese Gruppen
SupplementaryGroups=gpio audio
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
